# ha-energy-optimizer/rootfs/etc/services.d/energy-optimizer/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# HA Energy Optimizer - Service cleanup
# ==============================================================================
bashio::log.info "HA Energy Optimizer stopped."
